# bin/messenger/messenger.service.dist
[Unit]
Description=Symfony messenger-consume async
After=network.target

[Service]
Type=simple
User=apache
ExecStart=/var/www/wapinet/bin/console messenger:consume async --limit=1000 --memory-limit=128M --time-limit=3600 --ansi
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
