# bin/messenger/scheduler.service.dist
[Unit]
Description=Symfony messenger-consume scheduler
After=network.target

[Service]
Type=simple
User=apache
ExecStart=/var/www/wapinet/bin/console messenger:consume scheduler_default --limit=1000 --memory-limit=128M --time-limit=1500000 --ansi
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
